# _bootstrap/scripts/launchd/com.didio.token-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.didio.token-report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HUB_ROOT__/_bootstrap/scripts/token-report.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__HUB_ROOT__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HUB_ROOT__/logs/token-report.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HUB_ROOT__/logs/token-report.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
